# openvpn/vendor/openvpn.supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
directory=/etc/openvpn
command=/usr/sbin/openvpn --config /usr/share/doc/openvpn/samples/sample-config-files/loopback-client
user=root
autostart=true
autorestart=true
startsec=0
stopwaitsec=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
