# sndio/sndiod.initd
#!/sbin/openrc-run

name="sndiod"
pname="sndiod"
user="sndiod"
group="audio"
command="/usr/sbin/$pname"
command_args="$DAEMON_OPTS"
command_backdground="no"

pidfile="/run/$pname.pid"

depend() {
	need net
	after firewall
}
