# chrony/chrony-wait.service.d/10-use-unix-domain-sockets.conf
[Service]
User=chrony

ExecStart=
ExecStart=/usr/bin/chronyc -h /run/chrony/chronyd.sock waitsync 0 0.1 0.0 1

IPAddressAllow=
RestrictAddressFamilies=
RestrictAddressFamilies=AF_UNIX

ReadWritePaths=/run/chrony
UMask=0077
